# crontab/send_subscriber_email
0 17 * * 4 python /app/manage.py send_subscriber_email
